# dist/alertmanager/etc/systemd/system/alertmanager.service
[Unit]
Description=Prometheus Alertmanager
Documentation=https://github.com/prometheus/alertmanager

[Service]
ExecStart=/usr/bin/node_exporter
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
